# build/biominer-indexd.service
[Unit]
Description=BioMiner Indexd Daemon
After=network.target

[Service]
EnvironmentFile=/opt/biominer-indexd/biominer-indexd.conf
ExecStart=/opt/biominer-indexd/biominer-indexd -H 0.0.0.0 -p 3000 -c /opt/biominer-indexd/config.json -o --data-dir /opt/biominer-indexd/datasets
ExecReload=/bin/kill -HUP $MAINPID
Type=simple
KillMode=process
User=biominer
Group=biominer
# Restart=on-failure
# RestartSec=42s

[Install]
WantedBy=multi-user.target